# cleanup.cron
*/2 * * * * echo "Running cleanup..." >> /data/cron.log 2>&1
